# MyFirstVirus/trojan_penguin.service
[Unit]
Description = XXX

[Service]
ExecStart = /usr/bin/sh
Type = fork

[Install]
WanteBy = multy-user.target
